# templates/systemd-user/context-healthcheck.timer
[Unit]
Description=Run Unified Context Healthcheck every 5 minutes (user)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
