# templates/launchd/com.openviking.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.openviking.server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__SCRIPTS_DIR__/start_openviking.sh</string>
  </array>
  <key>WorkingDirectory</key><string>__SCRIPTS_DIR__</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>__LOG_DIR__/openviking_server_launchd.log</string>
  <key>StandardErrorPath</key><string>__LOG_DIR__/openviking_server_launchd.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key><string>__HOME__</string>
    <key>OPENVIKING_SKIP_PIP</key><string>1</string>
    <key>LITELLM_LOCAL_MODEL_COST_MAP</key><string>True</string>
  </dict>
</dict>
</plist>
